# units/zfs-getforce.service
[Unit]
Description=ZFS force detection

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/zfs-getforce
